# content/etc/init/upstart_xbmc_bridge.conf
#!upstart

start on started xbmc
stop on stopped xbmc

respawn

script
   python /usr/local/sbin/upstart_xbmc_bridge.py
end script
